# deploy/aphrnts-100-deploy.timer
[Unit]
Description=Run APHRNTs_100 auto-deploy check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
